# awesome.service
[Unit]
Description=Awesome Window Manager

[Service]
ExecStart=/usr/bin/awesome

[Install]
WantedBy=wm.target
