# blueman.service
[Unit]
Description=Blueman Applet

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=graphical.target

