# chromium.service
[Unit]
Description=Google Chrome Browser

[Service]
ExecStart=/usr/bin/chromium --restore-last-session

[Install]
WantedBy=wm.target
