# chromium.service
[Unit]
Description=Google Chrome Browser

[Service]
ExecStart=/usr/bin/chromium

[Install]
WantedBy=wm.target
