# ctags@.path
[Unit]
Description=Rebuild Ctags file when directory %I changes

[Path]
PathModified=/home/kstep/%I

[Install]
WantedBy=default.target
WantedBy=paths.target
