# dbeaver.service
[Unit]
Description=Database management GUI

[Service]
ExecStart=/usr/bin/dbeaver-ee

[Install]
WantedBy=graphical.target
