# deadbeef.service
[Unit]
Description=DeadBeef Music Player

[Service]
ExecStart=/usr/bin/deadbeef

[Install]
WantedBy=graphical.target
