# docked.target
[Unit]
Description=Units run on dock
Conflicts=undocked.target
