# emacs.service
[Unit]
Description=Emacs editor

[Service]
Type=simple
ExecStart=/usr/bin/emacs
Restart=always
RestartSec=3

[Install]
WantedBy=graphical.target
