# firefox.service
[Unit]
Description=Google Chrome Browser

[Service]
ExecStart=/usr/bin/firefox

[Install]
WantedBy=graphical.target
