# geary.service
[Unit]
Description=Geary e-mail client

[Service]
ExecStart=/usr/bin/geary

[Install]
WantedBy=graphical.target
