# google-chrome.service
[Unit]
Description=Google Chrome Browser

[Service]
ExecStart=/usr/bin/google-chrome-beta

[Install]
WantedBy=graphical.target
