# graphical.target
[Unit]
Description=User GUI
Requires=default.target
Requires=dbus.socket
After=default.target
After=xorg.target
AllowIsolate=yes

