# kbdd.service
[Unit]
Description=Very Simple Layout Switcher

[Service]
Restart=always
ExecStart=/usr/bin/kbdd --nodaemon

[Install]
WantedBy=wm.target
