# libreoffice.service
[Unit]
Description=LibreOffice

[Service]
ExecStart=/usr/bin/libreoffice

[Install]
WantedBy=graphical.target
