# musescore.service
[Unit]
Description=MuseScore notes editor

[Service]
ExecStart=/usr/bin/musescore

[Install]
WantedBy=graphical.target
