# nm-applet.service
[Unit]
Description=Network Manager Applet

[Service]
ExecStart=/usr/bin/nm-applet
Restart=always
RestartSec=0

[Install]
WantedBy=graphical.target
