# obs.service
[Unit]
Description=OBS Studio Live Streaming Client

[Service]
ExecStart=/usr/bin/obs

[Install]
WantedBy=graphical.target