# pasystray.service
[Unit]
Description=PulseAudio Volume Systray Control

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=graphical.target
