# popcorntime.service
[Unit]
Description=Popcorn Time Video Player

[Service]
ExecStart=/usr/bin/popcorntime

[Install]
WantedBy=graphical.target
