# reload-awesome.service
[Unit]
Description=Reload awesome WM

[Service]
ExecStart=/usr/bin/systemctl --user reload awesome.service
