# rescuetime.service
[Unit]
Description=RescueTime time tracker

[Service]
Type=simple
ExecStart=/usr/bin/rescuetime

[Install]
WantedBy=graphical.target
