# skype.service
[Unit]
Description=Skype

[Service]
ExecStartPre=-/usr/bin/rm -f ~/.Skype/shared.lck
ExecStart=/usr/bin/skype

[Install]
WantedBy=wm.target
