# skype.service
[Unit]
Description=Skype

[Service]
ExecStart=/usr/bin/skype

[Install]
WantedBy=default.target
