# slack.service
[Unit]
Description=Slack Messaging Client

[Service]
ExecStart=/usr/bin/slack

[Install]
WantedBy=graphical.target
