# ssh-agent.socket
[Unit]
Description=SSH Agent Socket

[Socket]
ListenStream=%t/ssh-agent.sock

[Install]
WantedBy=sockets.target
