# steam.service
[Unit]
Description=Steam Store Client

[Service]
ExecStart=/usr/bin/steam

[Install]
WantedBy=graphical.target
