# telegram.service
[Unit]
Description=Telegram Messaging Client

[Service]
ExecStart=/usr/bin/telegram-desktop

[Install]
WantedBy=graphical.target
