# thunderbird.service
[Unit]
Description=Thunderbird e-mail client

[Service]
ExecStart=/usr/bin/thunderbird

[Install]
WantedBy=graphical.target
