# tmux@.socket
[Unit]
Description=Tmux Socket (session %i)

[Socket]
ListenStream=/tmp/tmux-%U/%i

[Install]
WantedBy=sockets.target

