# unclutter.service
[Unit]
Description=Unclutter Mouse Hider

[Service]
ExecStart=/usr/bin/unclutter

[Install]
WantedBy=graphical.target
