# undocked.target
[Unit]
Description=Units run on undock
Conflicts=docked.target
