# viber.service
[Unit]
Description=Viber Messaging Client

[Service]
ExecStart=/usr/bin/viber

[Install]
WantedBy=graphical.target
