# vimb.service
[Unit]
Description=Vimb Browser

[Service]
ExecStart=/usr/bin/vimb
Restart=on-success

[Install]
WantedBy=wm.target
