# vkplayer.service
[Unit]
Description=VKontakte Music Player

[Service]
ExecStart=/usr/bin/vkplayer

[Install]
WantedBy=graphical.target
