# vlock.service
[Unit]
Description=Lock all consoles with vlock

[Service]
ExecStart=/usr/bin/vlock -a
