# vscode.service
[Unit]
Description=Visual Studio Code

[Service]
ExecStart=/usr/bin/code --wait

[Install]
WantedBy=graphical.target
