# wm.target
[Unit]
Description=Window Manager Started
Requires=xorg.target
After=xorg.target
AllowIsolate=yes

[Install]
WantedBy=default.target
