# xfce4.service
[Unit]
Description=XFCE4 Session

[Service]
ExecStart=/usr/bin/startxfce4

[Install]
WantedBy=default.target
Alias=xsession.service
