# xorg.target
[Unit]
Description=Xorg Server Started Target

