# xorg.target
[Unit]
Description=Xorg Server Started Target
Wants=window-manager.service
RefuseManualStart=yes
