# xorg@.socket
[Unit]
Description=Socket for Xorg at display %i

[Socket]
ListenStream=/tmp/.X11-unix/X%i

[Install]
WantedBy=default.target
